# catalogue/upload-server.service
[Unit]
Description=Catalogue image upload server
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/ajtpro/catalogue
Environment=UPLOAD_TOKEN=REPLACE_WITH_REAL_TOKEN
Environment=UPLOAD_PORT=8090
ExecStart=/usr/bin/python3 /home/ubuntu/ajtpro/catalogue/upload-server.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
